# .sync/com.webber.obsidian-main-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.webber.obsidian-main-sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/chenweisheng/Library/Mobile Documents/com~apple~CloudDocs/Notes-Work/.sync/sync_icloud_to_noteswork.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/chenweisheng/Library/Mobile Documents/com~apple~CloudDocs/Notes-Work/.sync/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chenweisheng/Library/Mobile Documents/com~apple~CloudDocs/Notes-Work/.sync/launchd.err.log</string>
  </dict>
</plist>
